# packaging/systemd/wireroom.service
[Unit]
Description=Wireroom local-first agent switchboard
[Service]
Type=simple
WorkingDirectory=%h/wireroom
EnvironmentFile=%h/.config/wireroom/env
UMask=0077
# harn:assume fresh-clone-install-proven-by-script ref=systemd-service
ExecStart=/usr/bin/node %h/wireroom/packages/cli/dist/index.js --data-dir %h/.wireroom up --static-root %h/wireroom/packages/web/dist --room desk --room-name Desk
# harn:end fresh-clone-install-proven-by-script
Restart=on-failure
RestartSec=5s
TimeoutStopSec=30s
KillMode=mixed

[Install]
WantedBy=default.target
